# etc/systemd/system/cf-dynamic.timer
[Unit]
Description=Timer for cf-dynamic

[Timer]
OnUnitActiveSec=10m
OnBootSec=10s

[Install]
WantedBy=timers.target